# debian/naspi-control.service
[Unit]
Description=NasPi Control

[Service]
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/naspi-control

[Install]
WantedBy=multi-user.target